# Linux/gfgs1.service
[Unit]
Description=CPU Stress command service

[Service]

ExecStart= dd if=/dev/zero of=/dev/null
MemoryLimit=201M
CPUShares=1024